# debian/aptcron-random-time.cron
SHELL=/bin/bash
* *     * * *   root    at -f /usr/bin/touch $(($RANDOM \% 24)):$(($RANDOM \% 60))
