# debian/aptcron-fixed-time.cron.d
# Run aptcron at 6:15 am.

# m h dom mon dow   user    command
6 15  * * *         root    aptcron
